# deploy/planner.service
[Unit]
Description=Planner — Socratic project intake server
Documentation=https://github.com/theturtleautomation/planner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=planner
Group=planner

# Paths
ExecStart=/usr/local/bin/planner-server --port 3100 --static-dir /opt/planner/web
WorkingDirectory=/opt/planner

# Environment
EnvironmentFile=-/etc/planner/planner.env

# LLM CLI isolation is handled by the Rust code (CliEnvironment in
# planner-core/src/llm/providers.rs). Each CLI invocation gets a fully
# clean environment via env_clear() + provider-specific env vars.
# The service PATH includes /opt/planner/bin/ so resolve_cli_binary()
# can find binaries both via direct path check and via the fallback `which`.
Environment=HOME=/opt/planner
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/opt/planner/bin

# Logging — goes to journalctl -u planner
StandardOutput=journal
StandardError=journal
SyslogIdentifier=planner

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/planner/data /opt/planner/cli-home /opt/planner/cli-sandbox /tmp/planner-worktrees
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes

# Resource limits
LimitNOFILE=65536

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
